# home/local/service/ssh-agent/run
#!/bin/sh

if [ -z "$USER" ]; then
	USER=$(whoami)
	export USER
fi

exec ssh-agent -D -a "/run/x-user/$USER/ssh-agent.socket"
