# home/local/service/i3lock/run
#!/bin/sh

[ -f ./conf ] && . ./conf

: ${DISPLAY:=:0}
export DISPLAY

exec xss-lock -- i3lock -n -i ~/Pictures/lock.png -c 1d2021
