# home/local/service/weechat/run
#!/bin/sh

tmux new -d -s weechat weechat
exec pause
