# home/local/service/weechat/run
#!/bin/sh

cd ../../../
: ${HOME:=$PWD}
export HOME

tmux new -d -s weechat weechat
exec pause
