# system/sv/runsvdir-danh/finish
#!/bin/sh

[ -f conf ] && . conf

: ${ZUSER:=danh}
ZHOME=$(eval echo ~$ZUSER)

sv -w600 force-stop "$ZHOME"/.local/service/*
sv exit "$ZHOME"/.local/service/*
