# home/local/service/ssoma/run
#!/bin/sh

exec snooze -H'*' -M/10 ssoma sync
