# home/local/service/ssoma/run
#!/bin/sh

exec snooze -H'*' -M3/10 -S15 ssoma sync
