# home/local/service/slstatus/run
#!/bin/sh

[ -f ./conf ] && . ./conf

: ${DISPLAY:=:0}
export DISPLAY

PATH="$HOME/.local/bin:/usr/local/bin:$PATH"

exec slstatus
