# system/etc/rc.local
#!/bin/sh

for f in /etc/zoneshot.d/start/*; do
	if [ -f "$f" ] && [ -x "$f" ]; then
		$f &
	fi
done

wait
