# sv/runsvdir-danh/run
#!/bin/sh

exec 2>&1

install -d -m 700 -o danh -g danh /run/user/danh
install -d -m 700 -o danh -g danh /run/user/danh/runit

exec chpst -udanh:danh runsvdir ~danh/.local/service
