# home/local/service/urxvtd/run
#!/bin/sh

cd ../../../
: ${HOME:=$PWD}
export HOME

PATH="$HOME/bin:$HOME/.local/bin:/usr/local/bin:$PATH"

exec urxvtd
