# system/etc/sv/runsvdir-sgn/finish
#!/bin/sh

[ -f conf ] && . conf

: ${ZUSER:=sgn}
ZHOME=$(eval echo ~$ZUSER)

sv -w600 force-stop "$ZHOME"/.local/service/*
sv exit "$ZHOME"/.local/service/*
