# system/etc/udev/rules.d/00-low-battery.rules
SUBSYSTEM=="power_supply", ATTR{type}=="Battery", ATTR{status}=="Discharging", ATTR{capacity}=="4", RUN+="/usr/bin/ZZZ"
