# home/local/service/mbsync/run
#!/bin/sh

exec snooze -H'*' -M/2 mbsync -a
