# home/local/service/mbsync/run
#!/bin/sh

HOME=$(eval echo "~$(id -un)")
export HOME

exec snooze -H'*' -M/2 mbsync -a
