# home/local/service/mbsync/run
#!/bin/sh

exec 2>/dev/null
exec snooze -H'*' -M/2 -S30 mbsync -a
