# home/local/service/battery/run
#!/bin/sh

thisdir="${0%/*}"

./battery.sh
exec sleep 180
