# home/local/service/battery/run
#!/bin/sh

thisdir="${0%/*}"

exec snooze -H'*' -M/3 "$thisdir/battery.sh"
